# beepy-poll-1.1/debian/beepy-idle.timer
[Unit]
Description=Beepy idle shutdown check

[Timer]
OnBootSec=1m
OnUnitActiveSec=1m
Unit=beepy-idle.service

[Install]
WantedBy=timers.target
