# beepy-poll-1.1/debian/beepy-idle.service
[Unit]
Description=Beepy idle shutdown check

[Service]
Type=oneshot
ExecStart=/sbin/beepy-idle
StandardOutput=journal
StandardError=journal
